# etc/systemd/system/v2ray.service
[Service]
ExecStart=/usr/bin/v2ray run -d "/etc/v2ray"

[Install]
WantedBy=default.target
